# panapitouch.service
[Unit]
Description=PanaPiTouch - Panasonic PTZ Camera Monitor
After=graphical.target network-online.target panapitouch-preload.service
Wants=graphical.target network-online.target
Requires=panapitouch-preload.service

[Service]
Type=simple
User=admin
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-0
Environment=QT_QPA_PLATFORM=wayland
Environment=DISPLAY=:0
WorkingDirectory=/home/admin/PanaPiTouch
ExecStart=/home/admin/PanaPiTouch/venv/bin/python /home/admin/PanaPiTouch/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target




